# scripts/pier.service
[Unit]
Description=Pier PaaS
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=pier
Group=docker
WorkingDirectory=/opt/pier
EnvironmentFile=/opt/pier/.env
ExecStart=/opt/pier/bin/pier
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/pier/data /opt/pier/.docker /tmp
ProtectHome=true
# Expose root's docker config (read-only) so `docker login` on the host
# transparently feeds Pier deploys. We mount it OUTSIDE /root because
# ProtectHome=true overwrites /root with an inaccessible tmpfs even when a
# bind targets /root/.docker. Leading `-` makes the source optional;
# install.sh creates /root/.docker so the bind always succeeds.
BindReadOnlyPaths=-/root/.docker:/opt/pier/host-docker
Environment=HOME=/opt/pier
Environment=DOCKER_CONFIG=/opt/pier/host-docker
Environment=GIT_CONFIG_NOSYSTEM=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pier

[Install]
WantedBy=multi-user.target
